# private_dot_config/systemd/user/agentmemory.service
[Unit]
Description=agentmemory server
After=network.target

[Service]
Type=simple
ExecStart=/home/mcrowe/.local/share/mise/installs/npm-varlock/latest/bin/varlock run -p %h/.agentmemory -- /home/mcrowe/.local/share/mise/installs/node/22/bin/agentmemory
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
